# settings/etc/supervisor/gunicorn.conf
[program:gunicorn]
command=gunicorn common.wsgi -c /vagrant/settings/etc/gunicorn.py