# settings/etc/supervisor/runserver.conf
[program:runserver]
directory=/vagrant
command=python manage.py runserver 0.0.0.0:8000
autostart=True
autorestart=True
